# launch/com.jarvis.imessage.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <!-- Service Label -->
    <key>Label</key>
    <string>com.jarvis.imessage</string>

    <!-- Program to run -->
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/node</string>
        <string>/Users/YOUR_USERNAME/Jarvis-v0/dist/cli/imessage-cli.js</string>
        <string>start</string>
    </array>

    <!-- Working directory -->
    <key>WorkingDirectory</key>
    <string>/Users/YOUR_USERNAME/Jarvis-v0</string>

    <!-- Environment variables -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>NODE_ENV</key>
        <string>production</string>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>

    <!-- Run at load -->
    <key>RunAtLoad</key>
    <true/>

    <!-- Keep alive (restart if crashes) -->
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
        <key>Crashed</key>
        <true/>
    </dict>

    <!-- Restart behavior -->
    <key>ThrottleInterval</key>
    <integer>60</integer>

    <!-- Logging -->
    <key>StandardOutPath</key>
    <string>/Users/YOUR_USERNAME/Jarvis-v0/logs/imessage-stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/YOUR_USERNAME/Jarvis-v0/logs/imessage-stderr.log</string>

    <!-- Nice level (lower priority) -->
    <key>Nice</key>
    <integer>1</integer>

    <!-- Process type -->
    <key>ProcessType</key>
    <string>Background</string>

    <!-- Only run when user is logged in -->
    <key>LimitLoadToSessionType</key>
    <string>Aqua</string>
</dict>
</plist>
